# auxiliary_files/config_files/mumble-server.service.in
[Unit]
Description=Mumble server
After=network.target

[Service]
Type=exec
ExecStart=/usr/bin/@MUMBLE_SERVER_BINARY_NAME@ -fg -ini /etc/mumble-server.ini
Restart=always
PrivateDevices=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
